# server/crawler-bot.service
[Unit]
Description=China Governance Crawler Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/china-governance
EnvironmentFile=/root/china-governance/.env
# .env must contain: TELEGRAM_BOT_TOKEN=..., TELEGRAM_USER_ID=..., and ANTHROPIC_API_KEY=...
Environment=PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/python3 /root/china-governance/server/crawler_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
