# server/daily-crawl.timer
[Unit]
Description=Run China Governance crawl pipeline daily at 6am UTC

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
